# autorun.service
[Unit]
Description=gpt-bot
After=multi-user.target

[Service]
User=student
Group=student
Type=simple
Restart=always
ExecStart=/usr/bin/python3 /home/student/gpt_bot/bot.py

[Install]
WantedBy=multi-user.target